# dist/fix-bonjour.timer
[Unit]
Description=Timer to fix Bonjour/mDNS after boot
After=network.target avahi-daemon.service

[Timer]
# Run 30 seconds after boot to ensure everything is up
OnBootSec=30
# Run once per day at 3 AM as maintenance
OnCalendar=03:00
AccuracySec=1min

[Install]
WantedBy=timers.target

